# packaging/update-vai.service
[Unit]
Description=Arch Update Vai (non-interactive)

[Service]
Type=oneshot
ExecStart=/usr/bin/update-vai --auto --yes
